# KinoZavqBot.conf
[program:KinoZavqBot]
command=/root/TG-BOTS/venv/bin/python -u /root/TG-BOTS/@KinoZavqBot/app.py
directory=/root/TG-BOTS/@KinoZavqBot
autostart=true
autorestart=true
user=root
stdout_logfile=/root/TG-BOTS/@KinoZavqBot/logs/KinoZavqBot.log
redirect_stderr=true
