# Tarjimakinolar777bot.conf
[program:Tarjimakinolar777bot]
command=/root/TG-BOTS/venv/bin/python -u /root/TG-BOTS/@Tarjimakinolar777bot/app.py
directory=/root/TG-BOTS/@Tarjimakinolar777bot
autostart=true
autorestart=true
user=root
stdout_logfile=/root/TG-BOTS/@Tarjimakinolar777bot/logs/Tarjimakinolar777bot.log
redirect_stderr=true
